# deploy/manager/llm-chat-manager.service
[Unit]
Description=llm-chat manager (Rust)
Documentation=https://github.com/Palakorn-Voramongkol/llm-chat
After=network-online.target zitadel.service
Wants=network-online.target

[Service]
Type=simple
# Run as the human `llm` user — needed because the worker spawns the
# claude CLI which reads ~/.claude state (credentials, history, projects).
# Replicating that under a system user requires more wiring than payoff.
User=llm
Group=llm
WorkingDirectory=/home/llm
EnvironmentFile=/etc/llm-chat/manager.env
# Put the user-local bin (where claude 2.1.x lives via npm/native installer)
# before the system one so the worker spawns the same claude an interactive
# shell would, and its first-run "PATH not set" banner doesn't leak into
# chat answers.
Environment=PATH=/home/llm/.local/bin:/usr/local/bin:/usr/bin
ExecStart=/usr/local/bin/llm-chat-manager
Restart=on-failure
RestartSec=5s
LimitNOFILE=65536

# Hardening — manager + worker need writable /tmp (xvfb temp dir) and
# the user's home (claude state, sqlite). PrivateTmp gives a private
# writable /tmp; ProtectHome=read-only would block the claude state.
ProtectSystem=strict
PrivateTmp=true
NoNewPrivileges=true
ReadWritePaths=/home/llm
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
